# scripts/qr-collector.user.service
[Unit]
Description=qr Bybit linear L2 + trades collector (user-level)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/qr
EnvironmentFile=%h/qr/.env
ExecStart=%h/.local/bin/uv run python collector.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
